# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
#
# shellcheck shell=bash
#
# Start nginx proxy service


set -e

bashio::log.info "Creating nginx.conf"
# shellcheck disable=SC2046

# shellcheck disable=SC2086
echo $JSON_CONF | tempio \
    -template /etc/nginx/nginx.conf.gtpl \
    -out /etc/nginx.conf

# start server
bashio::log.info "Running nginx..."
stat "/data/ssl/cert.pem -c %y" > /tmp/certificate_timestamp
exec nginx -c /etc/nginx.conf < /dev/null